# files/etc/init.d/auto-expand-overlay
#!/bin/sh /etc/rc.common

START=19

start() {
    # 检查是否已经扩展过
    [ -f "/etc/overlay-expanded" ] && return 0

    # 执行扩展脚本
    /usr/bin/auto-expand-overlay &
}

stop() {
    return 0
}
